# .config/systemd/user/sync-mail.service
[Unit]
Description=Synchronize IMAP4 and Maildir mailboxes

[Service]
Type=simple
ExecStart=%h/.local/bin/sync-mail
